# capslock/.config/autostart/capslock.desktop
[Desktop Entry]
Type=Application
Name=Capslock
Exec=/home/gruszecki/dotfiles/capslock.sh
