# capslock/.config/autostart/capslock.desktop
[Desktop Entry]
Type=Application
Name=Capslock
Exec=/home/gruszecki/bin/,f-caps
StartupNotify=false
Terminal=true
X-GNOME-Autostart-enabled=true
